# souspi/system/service/souspi-controller/run
#!/bin/sh

## make sure temptracker is running
sv start temptracker || exit 1

exec SousVideControllerApp
