# souspi/system/service/souspi-lcdui/run
#!/bin/sh

# make sure controller is running first
sv start souspi-controller || exit 1

exec SousVideLCDUI
